# setup/services/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin admin --skip-login --login-options "-f admin" --noissue %I $TERM