# setup/services/getty@tty1.service.d/override.conf
[Service]
PreExecStart=/bin/setupcon
ExecStart=
ExecStart=-/sbin/agetty --autologin admin --noissue --skip-login %I $TERM
