# setup/services/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty %I $TERM