# setup/services/getty@tty2.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin pi --noclear --noissue --skip-login %I $TERM